# packages/systemd/docker-plugin-vault.socket
[Unit]
Description=Docker plugin for Hashicorp Vault

[Socket]
ListenStream=/run/docker/plugins/vault.sock

[Install]
WantedBy=sockets.target
